# projects/arkhe/rootfs/etc/systemd/system/orbvm.target
[Unit]
Description=OrbVM Target
AllowIsolate=yes
